# deploy/alertran-sgd.service
[Unit]
Description=Alertran SGD API (uvicorn)
After=network.target

[Service]
User=root
WorkingDirectory=/opt/alertran_sgd
Environment="ROOT_PATH=/Alertran_SGD"
ExecStart=/opt/alertran_sgd/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8100 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
